# host_app/systemd/plant_receiver.service.example
[Unit]
Description=Plant MQTT Receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Plantproject-ESP32/host_app
ExecStart=/usr/bin/python3 /home/pi/Plantproject-ESP32/host_app/run_receiver.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
